# POC/SD/docker/keepalived-supervisord.conf
[supervisord]
nodaemon=true

[program:keepalived]
command=/usr/sbin/keepalived

[program:rsyslogd]
command=rsyslogd -n

[program:ssh]
command=service ssh start

[program:cron]
command=service cron start
