# SL/docker/keepalived-supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=service apache2 start

[program:rsyslogd]
command=rsyslogd -n
