# SA/docker/keepalived-supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=rsyslogd -n
